# install/skel/centos/6/etc/init/devpanel-bootstrap.conf
start on runlevel [345]

task

exec /usr/bin/devpanel bootstrap --yes
