# install/skel/common/etc/cron.d/devpanel-every-few-mins
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin

*/3 * * * *  root  cd / && run-parts /opt/webenabled/bin/packages/cron.d/cron.every-few-mins
