# bin/packages/cron.d/cron.hourly/s3-cleanup
#!/bin/bash

devpanel run cronjob cleanup s3 --yes
