# bin/packages/cron.d/cron.hourly/s3-uploads
#!/bin/bash

devpanel run cronjob uploads to s3 --yes
