# install/skel/ubuntu/14/etc/init/devpanel-bootstrap.conf
#
# devPanel bootstrap service

description	"devPanel bootstrap service"

start on started networking

task

script
/usr/bin/devpanel bootstrap --yes
end script
