# install/skel/centos/7/etc/systemd/system/devpanel-bootstrap.service
[Unit]
Description=bootstrap devPanel services
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/bin/devpanel bootstrap --yes
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
